# hostblock-update.service
[Unit]
Description=Update custom hosts file
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=hostblock -u
